# launchd/com.darwin.relay-v2.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd LaunchAgent for darwin-relay-v2 HTTP server.

  The relay binds :8090 and is the public-facing surface (via the
  cloudflare tunnel). If it dies the frontend's POST /v0/intents,
  GET /v0/intents/:cid, POST /v0/redeem, etc. all 502 — every
  in-flight user flow breaks visibly.

  Env vars MUST match what the worker reads — same
  DARWIN_RELAY_V2_RELAY_WALLET_HEX / CONTROLLER_HEX / BALI_* etc.
  See com.darwin.relay-v2-worker.plist; they need to stay in sync.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.darwin.relay-v2</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/eden/data/darwin/repos/darwin-relay/target/release/darwin_relay_v2</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>RUST_LOG</key>
        <string>info</string>
        <!-- 2026-07-04: HOME moved out of /tmp (macOS wipe killed the
             worker overnight). See worker plist header for context. -->
        <key>MIDEN_NETWORK</key>
        <string>testnet</string>
        <key>DARWIN_RELAY_V2_STORE</key>
        <string>/tmp/darwin-relay/relay-v2.sqlite</string>
        <key>DARWIN_RELAY_V2_MIDEN_STORE</key>
        <string>/Users/eden/data/darwin/.relay-miden-testnet/.miden/store.sqlite3</string>
        <key>DARWIN_RELAY_V2_MIDEN_KEYSTORE</key>
        <string>/Users/eden/data/darwin/.relay-miden-testnet/.miden/keystore</string>
        <key>DARWIN_RELAY_V2_RELAY_WALLET_HEX</key>
        <!-- 2026-07-04: fresh Public relay wallet on testnet v0.15,
             created in the permanent HOME. Same wallet as the worker. -->
        <string>0x66e7105ea36a7491325480accb7331</string>
        <key>DARWIN_RELAY_V2_CONTROLLER_HEX</key>
        <!-- 2026-06-23: v7 controller redeployed on testnet v0.15. -->
        <string>0x719bd3a14b42533115b1bcc8e02ea5</string>
        <key>DARWIN_RELAY_V2_OUTBOUND_MODE</key>
        <string>bridge_out_v1</string>
        <!-- Bali dormant; placeholder valid v0.15 hex so boot passes. -->
        <key>DARWIN_RELAY_V2_BALI_FAUCET_HEX</key>
        <string>0xb0411b0e0c4985115c03d034234110</string>
        <key>DARWIN_RELAY_V2_BALI_BRIDGE_HEX</key>
        <string>0xb0411b0e0c4985115c03d034234110</string>
        <key>DARWIN_RELAY_V2_BALI_BRIDGE_SERVICE</key>
        <string>https://miden-testnet-bridge.dev.eu-north-3.gateway.fm</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/tmp/darwin-relay/server.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/darwin-relay/server.log</string>
</dict>
</plist>
