# launchd/com.darwin.relay-v2-worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd LaunchAgent for darwin-relay-v2 background worker.

  Drives the on-chain legs the HTTP server cannot: drain inbound
  P2ID notes into the relay vault, emit atomic_deposit_note for
  POSITION_CREDITED intents, fire b2agg burns for REQUESTED
  redemptions, poll the Bali bridge service for matured claims and
  flip sepolia_release_tx.

  Env vars MUST mirror com.darwin.relay-v2.plist (same accounts,
  same backing sqlite/keystore). When updating one, update both.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.darwin.relay-v2-worker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/eden/data/darwin/repos/darwin-relay/target/release/darwin_relay_v2_worker</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>RUST_LOG</key>
        <string>info</string>
        <!-- 2026-06-25: testnet was migrated to Miden v0.15 on
             2026-06-23. Old store + keystore (v0.14 hexes, v0.14
             genesis) are now incompatible with the new chain. Point
             the worker at the fresh v0.15 home created during the
             redeploy. -->
        <key>MIDEN_NETWORK</key>
        <string>testnet</string>
        <key>DARWIN_RELAY_V2_STORE</key>
        <string>/tmp/darwin-relay/relay-v2.sqlite</string>
        <key>DARWIN_RELAY_V2_MIDEN_STORE</key>
        <string>/tmp/miden-relay-testnet-home/.miden/store.sqlite3</string>
        <key>DARWIN_RELAY_V2_MIDEN_KEYSTORE</key>
        <string>/tmp/miden-relay-testnet-home/.miden/keystore</string>
        <key>DARWIN_RELAY_V2_RELAY_WALLET_HEX</key>
        <!-- 2026-06-25: fresh Public relay wallet on testnet v0.15
             (the old wallet 0x6651e5ac… was on the dead v0.14 chain). -->
        <string>0x98494868d18218b12f09e9bef727a1</string>
        <key>DARWIN_RELAY_V2_CONTROLLER_HEX</key>
        <!-- 2026-06-23: v7 controller redeployed on testnet v0.15
             (was 0xbef7d2e89e9c3e006e10f959fa16d2 on v0.14). -->
        <string>0x719bd3a14b42533115b1bcc8e02ea5</string>
        <key>DARWIN_RELAY_V2_OUTBOUND_MODE</key>
        <string>bridge_out_v1</string>
        <!-- 2026-06-25: dETH faucet hex on testnet v0.15 (was
             0x0a7d175e… on v0.14). The Epoch bridge will mint dUSDC
             from its own faucet on Miden side; until Epoch publishes
             that hex, point at our local dETH faucet so worker has a
             concrete asset to scan/credit. -->
        <key>DARWIN_RELAY_V2_FAUCET_HEX</key>
        <string>0xb0411b0e0c4985115c03d034234110</string>
        <key>DARWIN_RELAY_V2_OUTBOUND_FAUCET_HEX</key>
        <string>0xb0411b0e0c4985115c03d034234110</string>
        <!-- 2026-06-25: native controller for Miden-native flow polling.
             Same as the v7 controller now on v0.15 testnet. Set to
             empty string to disable native polling. -->
        <key>DARWIN_RELAY_V2_NATIVE_CONTROLLER_HEX</key>
        <string>0x719bd3a14b42533115b1bcc8e02ea5</string>
        <!-- Bali path is dormant (outbound_mode=bridge_out_v1 + bridge
             returns 404 since the v0.15 migration). Worker still parses
             these at boot, so we point at any valid v0.15 testnet hex
             to keep the boot pass. The DETH faucet is reused as a
             placeholder — never actually called for Bali. -->
        <key>DARWIN_RELAY_V2_BALI_FAUCET_HEX</key>
        <string>0xb0411b0e0c4985115c03d034234110</string>
        <key>DARWIN_RELAY_V2_BALI_BRIDGE_HEX</key>
        <string>0xb0411b0e0c4985115c03d034234110</string>
        <key>DARWIN_RELAY_V2_BALI_BRIDGE_SERVICE</key>
        <string>https://miden-testnet-bridge.dev.eu-north-3.gateway.fm</string>
        <key>DARWIN_RELAY_V2_WORKER_INTERVAL_S</key>
        <string>30</string>
        <!-- 2026-06-10: Epoch dUSDC path. The Epoch bridge ships
             USD-pegged tokens to the Miden side, so the worker
             needs to value incoming amounts at $1 instead of the
             live ETH price (legacy 1Click dETH path). dUSDC on
             Miden is 6-decimal (verified via vault snapshot
             2026-06-10 16:52: 1e18 wei → 1e6 dUSDC base). Need
             wei_per_miden_base = 10^12 (18→6 conversion) and
             asset_decimals = 6. Unset to revert to dETH path. -->
        <key>DARWIN_RELAY_V2_ASSET_PRICE_USD</key>
        <string>1.0</string>
        <key>DARWIN_RELAY_V2_ASSET_DECIMALS</key>
        <string>6</string>
        <key>DARWIN_RELAY_V2_WEI_PER_MIDEN_BASE</key>
        <string>1000000000000</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/tmp/darwin-relay/worker.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/darwin-relay/worker.log</string>
</dict>
</plist>
