# launchd/com.darwin.relay-v2-worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd LaunchAgent for darwin-relay-v2 background worker.

  Drives the on-chain legs the HTTP server cannot: drain inbound
  P2ID notes into the relay vault, emit atomic_deposit_note for
  POSITION_CREDITED intents, fire b2agg burns for REQUESTED
  redemptions, poll the Bali bridge service for matured claims and
  flip sepolia_release_tx.

  Env vars MUST mirror com.darwin.relay-v2.plist (same accounts,
  same backing sqlite/keystore). When updating one, update both.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.darwin.relay-v2-worker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/eden/data/darwin/repos/darwin-relay/target/release/darwin_relay_v2_worker</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>RUST_LOG</key>
        <string>info</string>
        <key>DARWIN_RELAY_V2_STORE</key>
        <string>/tmp/darwin-relay/relay-v2.sqlite</string>
        <key>DARWIN_RELAY_V2_MIDEN_STORE</key>
        <string>/Users/eden/.miden/store.sqlite3</string>
        <key>DARWIN_RELAY_V2_MIDEN_KEYSTORE</key>
        <string>/Users/eden/.miden/keystore</string>
        <key>DARWIN_RELAY_V2_RELAY_WALLET_HEX</key>
        <!-- Old wallet 0xed3cd5be… was a PRIVATE account that drifted
             on 2026-06-04 (chain advanced to a commitment we never
             computed locally); miden-client provides no recovery path
             for private wallets. Replaced with PUBLIC wallet so future
             drift can be auto-recovered via Client::import_account_by_id. -->
        <string>0x6651e5ac8195b2000187593fff7e03</string>
        <key>DARWIN_RELAY_V2_CONTROLLER_HEX</key>
        <string>0xbef7d2e89e9c3e006e10f959fa16d2</string>
        <key>DARWIN_RELAY_V2_OUTBOUND_MODE</key>
        <!-- 2026-06-06: bridge_out_v1 now calls /demo/flows/outbound/
             submit (the mock's actual outbound path) instead of the
             inbound /v0/deposit/submit. b2agg is dormant — Bali
             gateway.fm still 404 — but kept as the fallback for when
             Bali revives. Switch back via `bali` value (see
             redemption_outbound_rot memory). -->
        <string>bridge_out_v1</string>
        <!-- 2026-06-06 17:15: Bridge mock was reset with a fresh
             MIDEN_MASTER_SEED_HEX after its solver wallet drifted into
             a private-account-irrecoverable state (same drift class as
             relay_wallet 0xed3cd5be… earlier today). New seed →
             deterministic new IDs. Captured from the post-reset first
             demo deposit's bridge log line "submitting Miden pay-to-id
             mint". Override BOTH the inbound-side FAUCET_HEX (worker
             vault snapshots + bridge_out V1 read) AND the outbound-side
             OUTBOUND_FAUCET_HEX (asset shipped back to user on redeem)
             so the worker sees the dETH the new mock actually mints
             instead of the now-stale 0xf45bad08… The relay wallet
             still holds 1.4e14 of the old faucet but the new mock
             doesn't recognise it. Each future bridge reseed re-rolls
             this — see mock_bridge_faucet_drift memory. -->
        <key>DARWIN_RELAY_V2_FAUCET_HEX</key>
        <string>0x0a7d175ed63ec5200fb2ced86f6aa5</string>
        <key>DARWIN_RELAY_V2_OUTBOUND_FAUCET_HEX</key>
        <string>0x0a7d175ed63ec5200fb2ced86f6aa5</string>
        <key>DARWIN_RELAY_V2_BALI_FAUCET_HEX</key>
        <string>0xe63ba7bc2c19ff603c52c67fa4426d</string>
        <key>DARWIN_RELAY_V2_BALI_BRIDGE_HEX</key>
        <string>0xc98bb07c188cd2500e13f68a069cdc</string>
        <key>DARWIN_RELAY_V2_BALI_BRIDGE_SERVICE</key>
        <string>https://miden-testnet-bridge.dev.eu-north-3.gateway.fm</string>
        <key>DARWIN_RELAY_V2_WORKER_INTERVAL_S</key>
        <string>30</string>
        <!-- 2026-06-10: Epoch dUSDC path. The Epoch bridge ships
             USD-pegged tokens to the Miden side, so the worker
             needs to value incoming amounts at $1 instead of the
             live ETH price (legacy 1Click dETH path). dUSDC on
             Miden is 6-decimal (verified via vault snapshot
             2026-06-10 16:52: 1e18 wei → 1e6 dUSDC base). Need
             wei_per_miden_base = 10^12 (18→6 conversion) and
             asset_decimals = 6. Unset to revert to dETH path. -->
        <key>DARWIN_RELAY_V2_ASSET_PRICE_USD</key>
        <string>1.0</string>
        <key>DARWIN_RELAY_V2_ASSET_DECIMALS</key>
        <string>6</string>
        <key>DARWIN_RELAY_V2_WEI_PER_MIDEN_BASE</key>
        <string>1000000000000</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/tmp/darwin-relay/worker.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/darwin-relay/worker.log</string>
</dict>
</plist>
