# launchd/com.darwin.cf-tunnel-relay.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd LaunchAgent for the Cloudflare Quick Tunnel that exposes
  darwin-relay-v2's :8090 HTTP server to a public *.trycloudflare.com
  URL. The tunnel wrapper script (cf-tunnel.sh) captures the issued
  URL and publishes it to /tmp/cf-relay-url.

  IMPORTANT: this is the account-less / quick tunnel. The URL CHANGES
  on every cloudflared restart (and Cloudflare may reclaim slugs at
  any moment). For a stable reviewer-facing URL switch to a named
  tunnel + cloudflare account; that's TODO.

  Install:
    cp launchd/com.darwin.cf-tunnel-relay.plist ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.darwin.cf-tunnel-relay.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.darwin.cf-tunnel-relay</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/eden/data/darwin/repos/darwin-infra/scripts/cf-tunnel.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>CF_TUNNEL_PORT</key>
        <string>8090</string>
        <key>CF_TUNNEL_URL_FILE</key>
        <string>/tmp/cf-relay-url</string>
        <key>CF_TUNNEL_LABEL</key>
        <string>cf-tunnel-relay</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/tmp/cf-relay.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/cf-relay.log</string>
</dict>
</plist>
