# launchd/com.darwin.bridge-solver-autoheal.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd LaunchAgent for bridge-solver-autoheal.sh

  The autoheal script wraps `docker compose logs -f` in an outer
  while-true that survives the EOF the follow stream emits when its
  watched container is restarted (the script's own primary action).
  That handles the in-script death cause. It does NOT handle external
  process termination: SIGKILL, OOM kill, macOS Activity Monitor force
  quit, accidental pkill, panic via set -e tripping on an unexpected
  edge in the reactor body. Empirically we have seen all of these.

  launchd's `KeepAlive` is the supervisor for that residual surface.
  If the autoheal exits for any reason, launchd respawns it after
  `ThrottleInterval` seconds (10s here — long enough to avoid a tight
  loop during boot-time docker daemon races, short enough that we
  don't miss an underflow window).

  Install:
    cp launchd/com.darwin.bridge-solver-autoheal.plist \
       ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) \
       ~/Library/LaunchAgents/com.darwin.bridge-solver-autoheal.plist

  Reload after edits:
    launchctl bootout gui/$(id -u)/com.darwin.bridge-solver-autoheal
    launchctl bootstrap gui/$(id -u) \
       ~/Library/LaunchAgents/com.darwin.bridge-solver-autoheal.plist

  Watch:
    tail -f /tmp/autoheal.log
    launchctl print gui/$(id -u)/com.darwin.bridge-solver-autoheal
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.darwin.bridge-solver-autoheal</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/eden/data/darwin/repos/darwin-infra/scripts/bridge-solver-autoheal.sh</string>
    </array>

    <!-- docker compose lives at /usr/local/bin/docker on this host;
         keep /opt/homebrew/bin + /usr/bin too so the script's other
         tools (date, sleep, etc.) resolve under either Intel or
         Apple-silicon homebrew layouts. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- Start at login + respawn on any exit. ThrottleInterval is
         the minimum gap between respawns; launchd applies its own
         exponential backoff on top if the agent keeps crash-looping. -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- Append to the same log path the manual nohup invocation
         used, so existing `tail -f /tmp/autoheal.log` muscle memory
         keeps working. -->
    <key>StandardOutPath</key>
    <string>/tmp/autoheal.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/autoheal.log</string>
</dict>
</plist>
