# launchd/com.darwin.bridge-mock-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Belt-and-suspenders watchdog for Brian's 1Click bridge mock
  (miden-testnet-bridge-sepolia-bridge-1). The container already has
  Docker's `unless-stopped` restart policy set, but Docker-for-Mac
  has been observed to skip restart on some exit codes. This watchdog
  re-`docker start`s the container if its State.Status != "running".

  Runs every 60 s. Silent on success. Restart events recorded in
  /tmp/bridge-watchdog.log.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.darwin.bridge-mock-watchdog</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/eden/data/darwin/repos/darwin-infra/scripts/bridge-mock-watchdog.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>

    <key>StartInterval</key>
    <integer>60</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/bridge-watchdog-plist.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/bridge-watchdog-plist.log</string>
</dict>
</plist>
